# luminari.service
[Unit]
Description=LuminariMUD Game Server
After=network.target mysql.service

[Service]
Type=forking
User=luminari
Group=luminari
WorkingDirectory=/home/luminari/luminarimud
ExecStart=/home/luminari/luminarimud/autorun.sh
ExecStop=/bin/touch /home/luminari/luminarimud/.killscript
Restart=always
RestartSec=60

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=4096
LimitCORE=infinity

[Install]
WantedBy=multi-user.target